# 5-appinfra/modules/htc-infra/kubernetes/storage/parallelstore/image/supervisord-template.conf
[supervisord]
nodaemon=true

[program:daos_agent]
command=/usr/bin/daos_agent -o /etc/daos/daos_agent.yml

[program:dfuse]
command=/usr/bin/dfuse %(ENV_DFUSE_ARGS)s
depends_on=daos_agent
